# supervisor.conf.d/collector.conf
[program:collector]
directory=/home/dc/Git/damien
command=python damien --dbname prod collector
user=dc
autostart=true
autorestart=true
