# supervisor.conf.d/daemon.conf
[program:daemon]
directory=/home/dc/Git/damien
command=python damien --dbname prod daemon
user=root
autostart=true
autorestart=true
