# supervisor.conf.d/httpd.conf
[program:httpd]
directory=/home/dc/Git/damien
command=python damien --dbname prod httpd --plugins ./series --cache_dir /tmp/httpd
user=dc
autostart=true
autorestart=true
