# supervisor.conf.d/mongo.tunnel.conf
[program:mongotunnel]
command=ssh -N -L27017:127.0.0.1:27017 dc@132.227.76.25
user=dc
autostart=true
autorestart=true
